# files/ubuntu-root/etc/s6-overlay/s6-rc.d/svc-xorg/run
#!/usr/bin/with-contenv bash

# bail early on wayland
if [[ "${PIXELFLUX_WAYLAND}" == "true" ]]; then
  sleep infinity
fi

# Cleanup
rm -f /tmp/.X1-lock

# Enable DRI3 support if detected
VFBCOMMAND=""
if ! which nvidia-smi && [ -e "/dev/dri/renderD128" ]; then
  VFBCOMMAND="-vfbdevice /dev/dri/renderD128"
fi
if [ ! -z ${DRINODE+x} ]; then
  VFBCOMMAND="-vfbdevice ${DRINODE}"
fi
if [ "${DISABLE_DRI3}" != "false" ]; then
  VFBCOMMAND=""
fi

# Calculate optimal Xvfb resolution based on user settings
# Priority:
#   1. MAX_RES (explicit override)
#   2. SELKIES_MANUAL_WIDTH/HEIGHT (legacy manual mode)
#   3. CUSTOM_RESOLUTION with calculated buffer space for dynamic resize
#   4. Default fallback (1920x1080)
DEFAULT_RES=""

if [ ! -z ${MAX_RES+x} ]; then
  # Explicit MAX_RES override (highest priority)
  DEFAULT_RES="${MAX_RES}"
  echo "[svc-xorg] Using explicit MAX_RES: ${DEFAULT_RES}"
elif [ -n "${SELKIES_MANUAL_HEIGHT}" ] || [ -n "${SELKIES_MANUAL_WIDTH}" ]; then
  # Legacy manual mode
  T_WIDTH="${SELKIES_MANUAL_WIDTH:-1024}"
  T_HEIGHT="${SELKIES_MANUAL_HEIGHT:-768}"
  if [ "${T_WIDTH}" = "0" ]; then T_WIDTH="1024"; fi
  if [ "${T_HEIGHT}" = "0" ]; then T_HEIGHT="768"; fi
  DEFAULT_RES="${T_WIDTH}x${T_HEIGHT}"
  echo "[svc-xorg] Using manual dimensions: ${DEFAULT_RES}"
elif [ -n "${CUSTOM_RESOLUTION}" ]; then
  # Calculate optimal Xvfb resolution from CUSTOM_RESOLUTION
  # Add buffer for dynamic resizing while keeping memory usage reasonable
  # XVFB_BUFFER_MULTIPLIER can be set via env (default: 2.0)
  # - 1.5x: Minimal overhead, supports up to 2880x1620 for 1920x1080 base
  # - 2.0x: Balanced (default), supports up to 3840x2160 (4K) for 1920x1080 base
  # - 3.0x: Maximum flexibility, supports up to 5760x3240 for 1920x1080 base
  CUSTOM_WIDTH="${CUSTOM_RESOLUTION%x*}"
  CUSTOM_HEIGHT="${CUSTOM_RESOLUTION#*x}"
  
  if [[ "${CUSTOM_WIDTH}" =~ ^[0-9]+$ ]] && [[ "${CUSTOM_HEIGHT}" =~ ^[0-9]+$ ]]; then
    # Use XVFB_BUFFER_MULTIPLIER or default to 2.0
    BUFFER_MULT="${XVFB_BUFFER_MULTIPLIER:-2.0}"
    
    # Validate buffer multiplier (1.0 - 4.0)
    if ! awk "BEGIN { exit !(${BUFFER_MULT} >= 1.0 && ${BUFFER_MULT} <= 4.0) }" 2>/dev/null; then
      echo "[svc-xorg] Warning: Invalid XVFB_BUFFER_MULTIPLIER ${BUFFER_MULT}, using default 2.0"
      BUFFER_MULT="2.0"
    fi
    
    # Calculate buffered resolution
    BUFFERED_WIDTH=$(awk "BEGIN {print int(${CUSTOM_WIDTH} * ${BUFFER_MULT})}")
    BUFFERED_HEIGHT=$(awk "BEGIN {print int(${CUSTOM_HEIGHT} * ${BUFFER_MULT})}")
    DEFAULT_RES="${BUFFERED_WIDTH}x${BUFFERED_HEIGHT}"
    echo "[svc-xorg] Calculated Xvfb resolution from CUSTOM_RESOLUTION ${CUSTOM_RESOLUTION}: ${DEFAULT_RES} (${BUFFER_MULT}x buffer)"
  else
    echo "[svc-xorg] Warning: Invalid CUSTOM_RESOLUTION format: ${CUSTOM_RESOLUTION}"
    DEFAULT_RES="1920x1080"
  fi
else
  # Fallback default (much more reasonable than 8K)
  DEFAULT_RES="1920x1080"
  echo "[svc-xorg] Using default resolution: ${DEFAULT_RES}"
fi

# Set DPI based on environment variable (default 96 for standard displays, use 144 or 192 for HiDPI)
DPI_VALUE="${DPI:-96}"

# Run Xvfb server with required extensions
TARGET_USER="${CUSTOM_USER:-${USER_NAME:-root}}"
exec s6-setuidgid "$TARGET_USER" \
  /usr/bin/Xvfb \
    "${DISPLAY}" \
    -screen 0 "${DEFAULT_RES}x24" \
    -dpi "${DPI_VALUE}" \
    +extension "COMPOSITE" \
    +extension "DAMAGE" \
    +extension "GLX" \
    +extension "RANDR" \
    +extension "RENDER" \
    +extension "MIT-SHM" \
    +extension "XFIXES" \
    +extension "XTEST" \
    +iglx \
    +render \
    -nolisten "tcp" \
    -ac \
    -noreset \
    -shmem \
    ${VFBCOMMAND}
